# src/local_cloud_agent/conf/local_cloud_agent.service
[Unit]
Description=Avahi mDNS/DNS-SD Stack

[Service]
ExecStart=/bin/sh

[Install]
WantedBy=multi-user.target

